# galeforce/conf/dropbear.conf
start on starting boot-services

task

script
  dropbear -R
end script