# galeforce/conf/dropbear.conf
start on starting boot-services

task

script
  /bin/dropbear -R
end script

