# galeforce/conf/dropbear.conf
start on starting boot-services

task

script
  /usr/local/galeforce/bin/dropbear -R
end script

